# supervisor/conf.d/langfuse.conf
; Langfuse — four supervised processes, each logging to stdout/stderr (Cloudron captures these).
; start.sh does setup (ownership, secrets, CLOUDRON_* -> Langfuse env) then exec's supervisord.
; ClickHouse + MinIO start immediately; the web/worker wrappers gate on readiness so migrations run
; before either Node process binds (field-guide gotcha 37).

[program:clickhouse]
command=/usr/bin/clickhouse-server --config-file=/etc/clickhouse-server/config.xml
user=cloudron
autostart=true
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:minio]
; MINIO_ROOT_USER / MINIO_ROOT_PASSWORD are exported by start.sh and inherited here.
; --address binds 0.0.0.0 so the Cloudron httpPorts proxy can reach MinIO on the container IP for the
; public 'blob' subdomain (presigned media). localhost still resolves here for server-side/event ops.
; The console stays localhost-only. Security on the public surface = S3 signatures + a private bucket.
command=/usr/bin/minio server /app/data/minio --address 0.0.0.0:9100 --console-address 127.0.0.1:9101
user=cloudron
autostart=true
autorestart=true
priority=10
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:web]
command=/app/code/conf/run-web.sh
directory=/app/code/web
user=cloudron
autostart=true
autorestart=true
priority=20
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=/app/code/conf/run-worker.sh
directory=/app/code/worker
user=cloudron
autostart=true
autorestart=true
priority=30
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
